# config/supervisor/bill.conf
[program:bill]
command=/usr/bin/python /home/lct/user.server/bill/app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/bill.err.log
stdout_logfile=/var/log/bill.out.log