# config/supervisor/charger.conf
[program:charge]
command=/usr/bin/python /home/lct/user.server/charger/app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/charge.err.log
stdout_logfile=/var/log/charge.out.log