# config/supervisor/charge.conf
[program:charge]
command=python3 /home/lct/user.server/charge/app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/charge.err.log
stdout_logfile=/var/log/charge.out.log