# config/supervisor/get_server.conf
[program:get_server]
command=/usr/bin/python /home/lct/user.server/get_server.py
autostart=true
autorestart=true
stderr_logfile=/var/log/get.err.log
stdout_logfile=/var/log/get.out.log