# config/supervisor/good.conf
[program:good]
command=python3 /home/lct/user.server/good/app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/good.err.log
stdout_logfile=/var/log/good.out.log