# config/supervisor/message.conf
[program:message]
command=python3 /home/lct/user.server/message/app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/message.err.log
stdout_logfile=/var/log/message.out.log