# config/supervisor/periodic.conf
[program:periodic]
command=python3 /home/lct/user.server/periodic.py
autostart=true
autorestart=true
stderr_logfile=/var/log/periodic.err.log
stdout_logfile=/var/log/periodic.out.log