# config/supervisor/recommend.conf
[program:recommend]
command=python3 /home/lct/user.server/recommend/app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/recommend.err.log
stdout_logfile=/var/log/recommend.out.log