# config/supervisor/user.conf
[program:user]
command=python3 /home/lct/user.server/user/app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/user.err.log
stdout_logfile=/var/log/user.out.log